# cookbooks/chef-caimito/templates/default/caimito.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:        caimito
# Required-Start:    
# Required-Stop:     
# Default-Start:   2 3 4 5
# Default-Stop:    
# Short-Description: Start caimito daemon
### END INIT INFO

# *****************************************************************
# Dynamic caimito configuration  --  DO NOT EDIT THIS FILE BY HAND
# Generated by chef - <%= @date %>
# *****************************************************************
PATH=/sbin:/bin:/usr/sbin:/usr/bin

. /lib/lsb/init-functions

do_start() {
  <%= @prefix %>/caimito/bin/demon.sh
}

do_stop()  {
  do_status && <%= @prefix %>/caimito/bin/shutdown.sh && sleep 2
  do_status && <%= @prefix %>/caimito/bin/shutdown.sh && sleep 5
  ! do_status || killall java
}

do_status(){
  ps ax|grep -s <%= @prefix %>/caimito/| grep -qsv grep
}

case "$1" in
start)
 do_start || (do_stop && false) ;;
restart|reload|force-reload)
 do_stop; do_start ;;
stop)
 do_stop ;;
status)
  if do_status; then echo "running..."; else  echo "stopped..."; false; fi;;
*)
 echo "Syntaxe: $(basename $0) [start|stop]..." && false
esac

