# cookbooks/chef-iscsi/templates/default/iscsi.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:        caimito
# Required-Start:    
# Required-Stop:     
# Default-Start:   2 3 4 5
# Default-Stop:    
# Short-Description: Start caimito daemon
### END INIT INFO

# *****************************************************************
# Dynamic iscsi configuration  --  DO NOT EDIT THIS FILE BY HAND
# Generated by chef - <%= @date %>
# *****************************************************************
PATH=/sbin:/bin:/usr/sbin:/usr/bin

. /lib/lsb/init-functions

do_start() {
  for j in <%= @portal.each do |ip|; ip; end -%>; do
    for i in $(iscsiadm --mode discovery --type sendtargets --portal $j| grep "$j:"| cut -d' ' -f2); do
      
      iscsiadm --mode node --targetname $i --portal $j --login
    done
  done
}

do_stop()  {
  true
}

do_status(){
  true
}

case "$1" in
start)
 do_start || (do_stop && false) ;;
restart|reload|force-reload)
 do_stop; do_start ;;
stop)
 do_stop ;;
status)
  if do_status; then echo "running..."; else  echo "stopped..."; false; fi;;
*)
 echo "Syntaxe: $(basename $0) [start|stop]..." && false
esac

